# local-cluster/etc/supervisor/conf.d/springtail-coordinator.conf
[program:springtail-coordinator]
directory=%(ENV_INSTALL_DIR)s/stc/coordinator
user=springtail
; Load env file, set unlimited core, then exec the app
command=/bin/sh -lc 'set -a; ulimit -c unlimited; exec %(ENV_INSTALL_DIR)s/venv/bin/python coordinator.py -c %(ENV_INSTALL_DIR)s/stc/coordinator/config.yaml --debug -f /opt/springtail/pids/postmaster.pid'

; Restart-on-failure (systemd’s Restart=on-failure)
autorestart=unexpected
exitcodes=0

; Environment from unit
environment=PYTHONUNBUFFERED="1",HOME="/home/springtail"

; disabled by default, enable with `supervisorctl enable springtail-coordinator`
autostart=false
startsecs=1
priority=50

; Signals & groups (closer to systemd semantics)
stopwaitsecs=35
stopsignal=TERM
stopasgroup=false
killasgroup=false

; Logging
stdout_logfile=/var/log/supervisor/springtail-coordinator.out.log
stderr_logfile=/var/log/supervisor/springtail-coordinator.err.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
stdout_logfile_backups=10
stderr_logfile_backups=10
